# systemd/omarchy-secureboot-watcher.service
[Unit]
Description=Repair Omarchy Secure Boot state after boot file changes

[Service]
Type=oneshot
ExecStart=@BINDIR@/omarchy-secureboot --quiet sign
